# crypto-bot.service
[Unit]
Description=v600-DYNAMO Dry Run (Production Mode)
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace/ai_crypto_strategy
# 建議在實機上建立 venv，或指向正確的 python 路徑
ExecStart=/usr/bin/python3 execution/live_dry_run.py
Restart=always
RestartSec=10
EnvironmentFile=/workspace/ai_crypto_strategy/.env

[Install]
WantedBy=multi-user.target
